# templates/ntp/ntp.service.erb
[Unit]
Description=NTP daemon
Before=multi-user.target
After=network-online.target
After=systemd-journald-dev-log.socket
Wants=network-online.target

[Service]
Type=forking
Restart=always
StartLimitInterval=20
StartLimitBurst=5
TimeoutStartSec=0
RestartSec=5
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
ExecStart=/etc/init.d/ntp start
ExecStop=/etc/init.d/ntp stop


[Install]
WantedBy=multi-user.target
