# templates/fticks/fticks.service.erb
[Unit]
Description=fticks
After=rsyslog.service

[Service]
ExecStart=/usr/sbin/fticks -f <%= @args %> -u <%= @url %> -p <%= @pipe %>

[Install]
WantedBy=multi-user.target
