# templates/fticks/fticks.service.erb
[Unit]
Description=fticks
After=rsyslog.service

[Service]
Type=simple
Restart=yes
StartLimitInterval=20
StartLimitBurst=5
TimeoutStartSec=0
RestartSec=5
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
ExecStart=/usr/sbin/fticks -d <%= @args %> -u <%= @url %> -p <%= @pipe %>


[Install]
WantedBy=multi-user.target
