# templates/unbound/sunet_unbound.service.erb
[Unit]
Description=Start unbound on SUNET host
After=network.target
Wants=docker.service
# We ought to add BindsTo if Docker is installed, but it
# makes unbound not start if Docker is NOT installed
#BindsTo=docker.service

[Service]
Restart=on-failure
StartLimitInterval=20
StartLimitBurst=5
TimeoutStartSec=0
ExecStart=/usr/sbin/unbound -d

[Install]
WantedBy=multi-user.target
Alias=unbound.service
