# templates/unbound/sunet_unbound.service.erb
[Unit]
Description=Start unbound on SUNET host
After=network.target
Wants=docker.target

[Service]
Restart=on-failure
StartLimitInterval=20
StartLimitBurst=5
TimeoutStartSec=0
ExecStart=/usr/sbin/unbound

[Install]
WantedBy=multi-user.target
Alias=unbound.service
