# templates/podmanhost/compose.service.erb
# This service file is generated by Puppet. Do not edit.
#
# This information is used by check_docker_containers to monitor the started containers:
# COMPOSE_FILE=<%= @compose_file %>
#
[Unit]
Description=<%= @description %>
Conflicts=shutdown.target reboot.target halt.target

[Service]
TimeoutStartSec=0
TimeoutStopSec=30
Restart=always
RestartSec=10

SyslogIdentifier=<%= @_service_name %>

WorkingDirectory=/

# To make docker pull work with basic-auth credentials
Environment="HOME=/root"

#ExecStartPre=-/usr/local/bin/podman-compose -f <%= @compose_file %> rm -f --stop
<% if @pull_on_start != true -%>#<% end -%>ExecStartPre=-/usr/local/bin/podman-compose -f <%= @compose_file %> pull

<% if @start_command -%>
ExecStart=<%= @start_command %>
<% else %>
ExecStart=/usr/local/bin/podman-compose -f <%= @compose_file %> up --force-recreate --no-deps
<% end -%>

ExecStop=/usr/local/bin/podman-compose -f <%= @compose_file %> down

NotifyAccess=all

<% if @service_extras.is_a? Array -%>
# service_extras :
<% @service_extras.each do |this| -%>
<%= this %>
<% end -%>
<% end -%>


[Install]
WantedBy=multi-user.target
