# templates/dockerhost/compose.service.erb
[Unit]
Description=<%= @description %>
After=docker.service
BindsTo=docker.service
Conflicts=shutdown.target reboot.target halt.target

[Service]
TimeoutStartSec=0
TimeoutStopSec=30
Restart=always
RestartSec=10

WorkingDirectory=/

ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> down
ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> kill
ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> rm
ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> rm -f

ExecStart=/usr/local/bin/docker-compose -f <%= @compose_file %> up --force-recreate --no-deps

ExecStop=/usr/local/bin/docker-compose -f <%= @compose_file %> stop

NotifyAccess=all

[Install]
WantedBy=multi-user.target
