# templates/dockerhost/compose.service.erb
# This service file is generated by Puppet. Do not edit.
#
# This information is used by check_docker_containers to monitor the started containers:
# COMPOSE_FILE=<%= @compose_file %>
#
[Unit]
Description=<%= @description %>
After=docker.service
BindsTo=docker.service
Conflicts=shutdown.target reboot.target halt.target

[Service]
TimeoutStartSec=0
TimeoutStopSec=30
Restart=always
RestartSec=10

WorkingDirectory=/

ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> down
ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> kill
ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> rm --all
ExecStartPre=-/usr/local/bin/docker-compose -f <%= @compose_file %> rm -f --all

ExecStart=/usr/local/bin/docker-compose -f <%= @compose_file %> up --force-recreate --no-deps

ExecStop=/usr/local/bin/docker-compose -f <%= @compose_file %> stop

NotifyAccess=all

[Install]
WantedBy=multi-user.target
